# templates/polkadot.service.j2
[Unit]
Description=Polkadot Dummy Node

[Service]
User=polkadot
Group=polkadot
ExecStart=/usr/local/bin/polkadot \
         --name polkadot-dummy \
         --validator \
         --rpc-methods=Unsafe \
         --chain={{ chain }} \
         --rpc-port {{ api_rpc_port }} \
         --database paritydb

[Install]
WantedBy=multi-user.target